# build/app.ini
[program:microservice]
user=user
directory=/var/www/html
command=/usr/local/bin/php ./bin/hyperf.php start
autostart=true
autorestart=true
stdout_logfile=/var/www/logs/app.log
stderr_logfile=/var/www/logs/errors.log
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0