# build/dev.ini
[program:microservice]
user=root
directory=/var/www/html
command=/usr/local/bin/php ./bin/hyperf.php server:watch
autostart=true
autorestart=true
stdout_logfile=/var/www/logs/app.log
stderr_logfile=/var/www/logs/errors.log
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0